# tp2bmc/board/tp2bmc/overlay/etc/profile.d/prompt.sh
if tty | grep -q '^/dev/pts'; then
	export PS1='\h \w \$ '
fi
